# com.getup.micmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label: unique identifier for this LaunchAgent -->
    <key>Label</key>
    <string>com.getup.micmonitor</string>

    <!-- Program: path to the mic monitor script -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>MIC_MONITOR_PATH_PLACEHOLDER</string>
    </array>

    <!-- Start automatically when the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if the script exits unexpectedly -->
    <key>KeepAlive</key>
    <true/>

    <!-- Throttle restarts: wait at least 10 seconds between restart attempts -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Log stdout and stderr -->
    <key>StandardOutPath</key>
    <string>LOG_PATH_PLACEHOLDER/mic-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_PATH_PLACEHOLDER/mic-monitor.log</string>
</dict>
</plist>
